# conf/api.conf
[program:api]
directory=/mnt/Detectron
command=/usr/bin/python api/api.py
stderr_logfile=/mnt/logs/api_logs/err.log
stdout_logfile=/mnt/logs/api_logs/out.log
autostart=true
autorestart=true